# root/etc/services.d/onlyoffice-example/run
#!/bin/sh
cd /usr/local/www/onlyoffice/example
exec /usr/local/bin/s6-setuidgid bsd /usr/local/bin/node ./bin/www
