# root/etc/services.d/rabbitmq/run
#!/bin/sh
exec /usr/local/sbin/rabbitmq-server
